# khal/vdirsyncer.service
[Service]
ExecStart=/home/bbguimaraes/.config/vdirsyncer/run.sh
Restart=always
