# isync/mbsync@.service
[Service]
ExecStart=/home/bbguimaraes/src/dotfiles/isync/mbsync.sh %i
Restart=always

[Install]
WantedBy=default.target
