# systemd/weather.service
[Service]
ExecStart=/home/bbguimaraes/src/dotfiles/scripts/weather.sh sync

[Install]
WantedBy=default.target
