# offlineimap/offlineimap@.service
[Service]
ExecStart=/usr/bin/offlineimap -u quiet -c ~/.config/offlineimap/%i
Restart=always

[Install]
WantedBy=default.target
